# crontab
PATH=/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin
LD_LIBRARY_PATH=/usr/local/lib

* * * * * ~/clash_royale_clan_data/clan_information_log.py >> ~/logs/clan_information_log.log 2>&1
